# apps/data/pgsync/supervisord.conf
[supervisord]
nodaemon=true

[program:pgsync]
command=bash -c 'bootstrap --config /app/config/schema.json && pgsync --config /app/config/schema.json --daemon'
autorestart=true
stopasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:dragonflydb]
command=dragonfly
autorestart=true
stopasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
